# OCodeControllerBackend/ServiceDockerFile/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

[program:opencode]
# OpenCode 智能体服务
command=opencode web --hostname 0.0.0.0 --port 4096
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:filebrowser]
# FileBrowser 文件管理服务
command=/usr/local/bin/filebrowser --port 8080 --address 0.0.0.0 --root / --database /var/lib/filebrowser/filebrowser.db
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
